# zero-deploy/deploy-script/script-linux/service/model-job-PvInerter.service
[Unit]
Description=model-job-PvInerter server
#Documentation=none
After=network.target

[Service]
Environment="SERVER_HOME={{.app.home}}/model-job-PvInerter"
ExecStart=/usr/local/bin/java -jar ${SERVER_HOME}/model-job-PvInerter.jar --spring.config.location=file:${SERVER_HOME}/application.yml
ExecStop=/bin/kill -9 $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
Alias=model-job-PvInerter.service

